# kotlin-source/src/main/resources/nodeConfiguration/notary.service
[Unit]
Description=Corda Node - Notary
Requires=network.target

[Service]
Type=simple
