# kotlin-source/src/main/resources/nodeConfiguration/partyb-webserver.service
[Unit]
Description=Webserver for Corda Node - Party B
Requires=network.tar