# kotlin-source/src/main/resources/nodeConfiguration/partyc-webserver.service
[Unit]
Description=Webserver for Corda Node - P